# website/deploy/cryptex-website.service.example
[Unit]
Description=CryptEX Website
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/cryptex/website
Environment=CRYPTEX_SITE_URL=https://cryptexorg.duckdns.org
Environment=CRYPTEX_ALLOWED_HOSTS=cryptexorg.duckdns.org,localhost,127.0.0.1
Environment=CRYPTEX_FORCE_HTTPS=true
Environment=CRYPTEX_CANONICAL_REDIRECT=true
Environment=CRYPTEX_SECURITY_CONTACT=Anon-Sec-BTCC@proton.me
Environment=CRYPTEX_REPOSITORY_URL=https://github.com/Anonymous137-sudo/CryptEX_Core
Environment=CRYPTEX_SITE_HOST=127.0.0.1
Environment=CRYPTEX_SITE_PORT=8080
Environment=CRYPTEX_SITE_RELOAD=0
ExecStart=/srv/cryptex/website/.venv/bin/uvicorn app.main:app --app-dir /srv/cryptex/website --host 127.0.0.1 --port 8080 --no-server-header --no-date-header
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
